# tags/pardus-1.0/system/base/hotplug/files/coldplug
#!/sbin/runscript
# Copyright © 2005  TUBITAK/UEKAE
# Licensed under the GNU General Public License, version 2.
# See the file http://www.gnu.org/copyleft/gpl.txt.
#
# Original work belongs Gentoo Linux

depend() {
	need modules
}

checkconfig() {
	if [ ! -d /etc/hotplug ] ; then
		eerror "Coldplug requires scripts & configs in /etc/hotplug !"
		return 1
	fi
}

start() {
	checkconfig || return 1
	
	for RC in /etc/hotplug/*.rc
	do
		NAME=`basename $RC .rc`
		ebegin "Coldplugging $NAME devices"
		# We do not want to check the return status, as
		# some of the scripts may fail due to drivers not
		# compiled as modules ...
		$RC start
		eend 0
	done

	# Ugly hack to make sure we have ppp node, as udev may not catch up with pppd
	if [ ! -c /dev/ppp ] ; then
		/usr/bin/mknod -m 660 /dev/ppp c 108 0
		/usr/bin/chown root:dialout /dev/ppp
	fi
}

